# services/core/run
#!/usr/bin/with-contenv sh


/wait-for-postgres.sh -- /init-database.sh
cd /harbor/
envfile /etc/core/env /usr/local/bin/harbor_core