# services/clair-adapter/run
#!/usr/bin/with-contenv sh

exec /wait-for-postgres.sh -- envfile /etc/clair-adapter/env /usr/local/bin/clair-adapter