# services/clair/run
#!/usr/bin/with-contenv sh

/wait-for-postgres.sh -- /usr/local/bin/clair -config /etc/clair/config.yml