# services/clair/run
#!/usr/bin/with-contenv sh

/usr/local/bin/clair -config /etc/clair/config.yml