# services/jobservice/run
#!/usr/bin/with-contenv sh

/wait-for-postgres.sh -- envfile /etc/jobservice/env /usr/local/bin/harbor_jobservice -c /etc/jobservice/config.yml