# services/registryctl/run
#!/usr/bin/with-contenv sh

exec /wait-for-postgres.sh -- envfile /etc/registryctl/env /usr/local/bin/harbor_registryctl -c /etc/registryctl/config.yml