# services/chartmuseum/run
#!/usr/bin/with-contenv sh

/wait-for-postgres.sh -- envfile /etc/chartmuseum/env /usr/local/bin/chartm