# services/registry/run
#!/usr/bin/with-contenv sh
exec /wait-for-postgres.sh -- /usr/local/bin/registry serve /etc/registry/config.yml