# services/redis/run
#!/usr/bin/with-contenv sh
/usr/local/bin/redis-server /etc/redis.conf