# services/portal/run
#!/usr/bin/with-contenv sh
exec /usr/sbin/nginx -c /etc/portal/nginx.conf -g "daemon off;"